# debian-jessie/authentic2-multitenant.cron.hourly
#!/bin/sh

authentic2-multitenant-manage tenant_command --all clearsessions
authentic2-multitenant-manage tenant_command --all cleanupauthentic
authentic2-multitenant-manage tenant_command --all sync-ldap-users
