# authentic2.cron.hourly
#!/bin/sh

/usr/sbin/service authentic2 manage clearsessions
/usr/sbin/service authentic2 manage cleanupauthentic
/usr/sbin/service authentic2 manage sync-ldap-users
